# systemd/nekonet-autoupdate.service
[Unit]
Description=NekoNet AutoUpdate v0.2 Coordinator
Wants=network-online.target
After=network-online.target
[Service]
Type=simple
EnvironmentFile=/etc/nekonet-autoupdate.env
ExecStart=/opt/nekonet-autoupdate/.venv/bin/nekonet-autoupdate
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
