# systemd/nekonet-autoupdate-run.service
[Unit]
Description=Start NekoNet maintenance run
After=nekonet-autoupdate.service
[Service]
Type=oneshot
EnvironmentFile=/etc/nekonet-autoupdate.env
ExecStart=/usr/bin/curl -fsS -X POST -H "Authorization: Bearer ${NEKONET_API_TOKEN}" http://${NEKONET_API_BIND}:${NEKONET_API_PORT}/api/v1/runs/start
