# systemd/nekonet-autoupdate.service
[Unit]
Description=NekoNet AutoUpdate v0.1 Coordinator
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/nekonet-autoupdate.env
ExecStart=/opt/nekonet-autoupdate/.venv/bin/nekonet-autoupdate
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nekonet-autoupdate

[Install]
WantedBy=multi-user.target
